# files/etc/profile.d/eternal-just.sh
alias just="just --unstable --justfile /usr/share/eternal/main.just"
