# files/usr/etc/profile.d/eternal-just.sh
alias eternal="just --justfile /usr/share/eternal/justfile"
